# systemd/company-os-sync.service
[Unit]
Description=Company OS S3 Sync Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/ai-meeting-notes-agent/scripts/sync-all.sh
Restart=always
RestartSec=10
Environment=SYNC_INTERVAL=300
Environment=PATH=%h/.local/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
